# mvidarr.service
[Unit]
Description=MVidarr Enhanced - Professional Music Video Management System with Advanced Processing
Documentation=file:///home/mike/mvidarr/README.md file:///home/mike/mvidarr/MILESTONE_ROADMAP.md file:///home/mike/mvidarr/PHASE_2_WEEKS18-24_ROADMAP.md
After=network-online.target mysql.service mariadb.service
Wants=network-online.target
RequiresMountsFor=/home/mike/mvidarr

[Service]
Type=exec
User=mike
Group=mike
WorkingDirectory=/home/mike/mvidarr

# Create required directories and start application
ExecStartPre=/usr/bin/mkdir -p /home/mike/mvidarr/data/logs /home/mike/mvidarr/data/downloads /home/mike/mvidarr/data/thumbnails /home/mike/mvidarr/data/cache /home/mike/mvidarr/data/backups /home/mike/mvidarr/data/processing
ExecStart=/home/mike/mvidarr/venv/bin/python /home/mike/mvidarr/app_launcher_v2.py fastapi_hybrid
ExecReload=/bin/kill -HUP $MAINPID

# Enhanced restart and failure handling
Restart=always
RestartSec=20
StartLimitInterval=600
StartLimitBurst=3
TimeoutStartSec=120
TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGTERM

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mvidarr

# Environment variables for Phase 2 FastAPI Application
Environment=PATH=/home/mike/mvidarr/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/mike/mvidarr
Environment=FLASK_ENV=production
Environment=FASTAPI_ENV=production
Environment=PYTHONUNBUFFERED=1
Environment=ASYNC_MODE=enabled
Environment=PHASE_2_ADVANCED_PROCESSING=enabled
Environment=FFMPEG_ADVANCED_TASKS=enabled
Environment=CONCURRENT_VIDEO_PROCESSING=enabled
Environment=QUALITY_ANALYSIS_ENABLED=enabled
Environment=BULK_OPERATIONS_ENABLED=enabled

# Security settings for home directory access
NoNewPrivileges=yes
PrivateTmp=yes
# Commented out for home directory access
# ProtectSystem=strict
# ProtectHome=yes
# ReadWritePaths=/home/mike/mvidarr/data
# ReadOnlyPaths=/home/mike/mvidarr

# Enhanced resource limits for Phase 2 advanced processing
LimitNOFILE=131072
LimitNPROC=8192

[Install]
WantedBy=multi-user.target